# services.d/teaspeak/run
#!/usr/bin/execlineb -P
s6-setuidgid teaspeak
/opt/teaspeak/teastart_minimal.sh