# root/etc/services.d/teaspeak/run
#!/usr/bin/execlineb -P
/opt/teaspeak/teastart_minimal.sh